# plists/com.skimthecream.email-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.skimthecream.email-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/erikbowen/erdabo-labs/skimthecream/scripts/run-service.sh</string>
    <string>email-watcher</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/erikbowen/erdabo-labs/skimthecream</string>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/erikbowen/erdabo-labs/skimthecream/logs/email-watcher.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/erikbowen/erdabo-labs/skimthecream/logs/email-watcher-error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
